# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgresql/data -c listen_addresses="*"
user=postgres
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/postgres.log
environment=POSTGRES_INITDB_ARGS="-U postgres"

[program:redis]
command=/usr/bin/redis-server --port 6379 --bind 0.0.0.0 --appendonly yes --dir /data
user=redis
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/redis.log

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
